# deploy/systemd/bmw-hai-dlab-preview.service
[Unit]
Description=BMW HAI D. LAB — Vite preview (API + static)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/bmw-hai-dlab
Environment=NODE_ENV=production
# Vite loads .env from WorkingDirectory via vite.config.ts
ExecStart=/opt/bmw-hai-dlab/node_modules/.bin/vite preview --host 0.0.0.0 --port 4173
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
